# udevrules/80-hotplug-cpu-mem.rules
# Do not edit this file under /usr/lib/udev/rules.d, it will be overwritten on update
# If you need modifications, copy it to /etc/udev/rules.d and edit it there

ACTION!="add", GOTO="hotplug_cpu_end"
SUBSYSTEM!="cpu", GOTO="hotplug_cpu_end"

#
# Hotplug physical CPU
#
CONST{arch}=="x86-64", TEST=="online", ATTR{online}=="0", ATTR{online}="1"
CONST{arch}=="s390*", TEST=="online", ATTR{configure}=="1", ATTR{online}=="0", ATTR{online}="1"
LABEL="hotplug_cpu_end"

#
# Hotplug physical memory
#
SUBSYSTEM!="memory", GOTO="hotplug_memory_end"
# Automatically set new memory devices online
# This causes a follow-up uevent with ACTION=="online"
ACTION=="add", CONST{arch}!="s390x", ATTR{state}=="offline", ATTR{state}="online"

# When memory size increases because memory device become online, instances of
# tmpfs are remounted so that their size is recalculated (boo#869603).
# This only affects tmpfs instances that are mounted with an explicit "size=X%"
# option in /etc/fstab.
ACTION=="online", RUN+="/usr/bin/systemctl --no-block start remount-tmpfs.service"
LABEL="hotplug_memory_end"
