# remount-tmpfs.service
[Unit]
Description=Remount tmpfs file systems

[Service]
Type=oneshot
ExecStart=/usr/lib/udev/remount-tmpfs
